# LUB/crontab
#! /bin/bash

30 11 * * 0 /home/blox_land/LUB/LUB.sh
00 11 * * * /home/blox_land/LUB/local_startupscript.sh

# $ crontab -e
